# supervisord.conf
[supervisord]
loglevel=info

[program:subscriptions]
directory=/var/www/project
# Too many docker layers, priviledged supervisord can't `run` new container because files are not visible/mountable, so I `exec`ed to already running one.
#command=docker-compose run --rm --no-deps php bin/console streak:subscriptions:run
command=docker-compose exec -T fpm php bin/console streak:subscriptions:run
autostart=true
autorestart=true
startsecs=0
startretries=1000
logfile = /var/www/project/var/log/supervisord.log
stdout_logfile = /var/www/project/var/log/supervisord.log
stderr_logfile = /var/www/project/var/log/supervisord.log
